# rm_old_files.cron
*/5 * * * * /home/pi/BSidesDCReg/rm_old_files.sh
