# rm_old_files.cron
*/5 * * * * /home/pi/PiConReg/rm_old_files.sh
